# supervisor/config/testdata/test.conf
[program:drpc]
command = /opt/drpc --env=test
process_name=%(program_name)s_%(process_num)d
numprocs=3
stdout_logfile = /opt/drpc/stdout.log
stderr_logfile = /opt/drpc/stderr.log
user = user_00
autostart = 1
autorestart = 1
startsecs = 10
stopasgroup = 1
killasgroup = 1
stdout_logfile_maxbytes = 500MB
stdout_logfile_backups = 10
stderr_logfile_maxbytes = 500MB
stderr_logfile_backups = 10

[program:dmicro]
command = /opt/dmicro --env=test
process_name=%(program_name)s_%(process_num)d
numprocs=3
stdout_logfile = /opt/dmicro/stdout.log
stderr_logfile = /opt/dmicro/stderr.log
user = user_00
autostart = 1
autorestart = 1
startsecs = 10
stopasgroup = 1
killasgroup = 1
stdout_logfile_maxbytes = 500MB
stdout_logfile_backups = 10
stderr_logfile_maxbytes = 500MB
stderr_logfile_backups = 10